# jacred/Data/crontab
*/5  *   *   *   *    curl -s "http://127.0.0.1:9118/jsondb/save"

*/2  *   *   *   *    curl -s "http://127.0.0.1:9118/cron/torrents/UpdateSizeInfo"
1    1   *   *   *    curl -s "http://127.0.0.1:9118/cron/torrents/UpdateSizeInfo?forced=true"
*/2  *   *   *   *    curl -s "http://127.0.0.1:9118/cron/torrents/UpdateDetails"
1    1   *   *   *    curl -s "http://127.0.0.1:9118/cron/torrents/UpdateDetails?forced=true"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9118/cron/torrents/UpdateTorrMediaCache"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/underverse/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9118/cron/underverse/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9118/cron/underverse/ParseAllTask"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/torrentby/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9118/cron/torrentby/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9118/cron/torrentby/ParseAllTask"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/rutracker/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9118/cron/rutracker/UpdateTasksParse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9118/cron/rutracker/ParseAllTask"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9118/cron/rutracker/parseMagnet"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/rutor/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9118/cron/rutor/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9118/cron/rutor/ParseAllTask"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/nnmclub/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9118/cron/nnmclub/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9118/cron/nnmclub/ParseAllTask"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/kinozal/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9118/cron/kinozal/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9118/cron/kinozal/ParseAllTask"
*    */4 *   *   *    curl -s "http://127.0.0.1:9118/cron/kinozal/parseMagnet"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/bitru/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9118/cron/bitru/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9118/cron/bitru/ParseAllTask"
*    */4 *   *   *    curl -s "http://127.0.0.1:9118/cron/bitru/parseMagnet"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/selezen/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9118/cron/selezen/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9118/cron/selezen/ParseAllTask"

*/40 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/toloka/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9118/cron/toloka/parseMagnet"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/animelayer/parse"
*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/animedia/parse"
*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/anilibria/parse"
*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/anifilm/parse"
*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/anidub/parse"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/baibako/parse"
*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/hamsterstudio/parse"
*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/hdrezka/parse"
*/15 *   *   *   *    curl -s "http://127.0.0.1:9118/cron/lostfilm/parse"
